# etc/supervisor/conf.d/dashboard_slave.conf
[program:dashboard_sync_files]
command=dashboard_sync_files.py
user=pi
directory=/home/pi/
autostart=true
autorestart=false
environment=PYTHONUNBUFFERED=1
