# etc/supervisor/conf.d/dashboard_master.conf
[program:dashboard_redis_import]
command=dashboard_redis_import.py
user=pi
directory=/home/pi/
autostart=true
autorestart=false
environment=PYTHONUNBUFFERED="1"

[program:dashboard_build_files]
command=dashboard_build_files.py
user=pi
directory=/home/pi/
autostart=true
autorestart=false
environment=PYTHONUNBUFFERED="1"

[program:dashboard_hmi]
command=dashboard_hmi.py
user=pi
directory=/home/pi/
autostart=true
autorestart=false
environment=PYTHONUNBUFFERED="1", DISPLAY=":0.0"
