# etc/supervisor/conf.d/dashboard_master.conf
[program:dashboard_redis_import]
command=dashboard_redis_import.py
user=pi
directory=/home/pi/
autostart=true
autorestart=false
environment=PYTHONUNBUFFERED=1
