# etc/supervisor/conf.d/dashboard_slave_loos.conf
[program:board_sync_files]
command=/opt/tk-dashboard/bin/board-sync-files.py
user=pi
directory=/home/pi/
autostart=true
autorestart=true
environment=PYTHONUNBUFFERED="1"

[program:board_hmi_loos]
command=/opt/tk-dashboard/bin/board-hmi-loos.py
user=pi
directory=/home/pi/
autostart=true
autorestart=true
environment=PYTHONUNBUFFERED="1", DISPLAY=:0.0
