# etc/supervisor/conf.d/dashboard_master_loos.conf
[program:dashboard_redis_import_loos]
command=dashboard_redis_import_loos.py
user=pi
directory=/home/pi/
autostart=true
autorestart=true
environment=PYTHONUNBUFFERED="1"

[program:dashboard_redis_export_loos]
command=dashboard_redis_export_loos.py
user=pi
directory=/home/pi/
autostart=true
autorestart=true
environment=PYTHONUNBUFFERED="1"

[program:dashboard_build_files_loos]
command=dashboard_build_files_loos.py
user=pi
directory=/home/pi/
autostart=true
autorestart=true
environment=PYTHONUNBUFFERED="1"

[program:dashboard_files_manager_loos]
command=dashboard_files_manager_loos.py
user=pi
directory=/home/pi/
autostart=true
autorestart=true
environment=PYTHONUNBUFFERED="1"

[program:dashboard_meters_manager_loos]
command=dashboard_meters_manager_loos.py
user=pi
directory=/home/pi/
autostart=true
autorestart=true
environment=PYTHONUNBUFFERED="1"

[program:dashboard_hmi_loos]
command=dashboard_hmi_loos.py
user=pi
directory=/home/pi/
autostart=true
autorestart=true
environment=PYTHONUNBUFFERED="1", DISPLAY=":0.0"
