# etc/supervisor/conf.d/dashboard_slave_loos.conf
[program:dashboard_sync_files_loos]
command=dashboard_sync_files_loos.py
user=pi
directory=/home/pi/
autostart=true
autorestart=true
environment=PYTHONUNBUFFERED="1"

[program:dashboard_hmi_loos]
command=dashboard_hmi_loos.py
user=pi
directory=/home/pi/
autostart=true
autorestart=true
environment=PYTHONUNBUFFERED="1", DISPLAY=":0.0"