# etc/supervisor/conf.d/dashboard_slave.conf
[program:dashboard_sync_files]
command=dashboard_sync_files.py
user=pi
directory=/home/pi/
autostart=true
autorestart=false
environment=PYTHONUNBUFFERED="1"

[program:dashboard_hmi]
command=dashboard_hmi.py
user=pi
directory=/home/pi/
autostart=true
autorestart=false
environment=PYTHONUNBUFFERED="1", DISPLAY=":0.0"