# etc/supervisor/conf.d/dashboard_slave_loos.conf
[program:dashboard_redis_ssh_forward]
command=autossh -M 0 -NT -o ServerAliveInterval=10 -o ServerAliveCountMax=3 -L 63790:127.0.0.1:6379 pi@192.168.0.60
priority=1
user=pi
directory=/home/pi/
autostart=true
autorestart=true

[program:dashboard_sync_files_loos]
command=dashboard_sync_files_loos.py
priority=2
user=pi
directory=/home/pi/
autostart=true
autorestart=true
environment=PYTHONUNBUFFERED="1"

[program:dashboard_hmi_loos]
command=dashboard_hmi_loos.py
priority=3
user=pi
directory=/home/pi/
autostart=true
autorestart=true
environment=PYTHONUNBUFFERED="1", DISPLAY=:0.0
